# data/export/upstart/process.conf.erb
start on starting <%= app %>-<%= name %>
stop on stopping <%= app %>-<%= name %>
respawn

env PORT=<%= port %>
<% engine.env.each do |name,value| -%>
env <%= name.upcase %>='<%= value.gsub(/'/, "'\"'\"'") %>'
<% end -%>

setuid <%= user %>

chdir <%= engine.root %>

#exec chruby-exec $FOREMAN_RUBY_VERSION -- <%= process.command %>

script
exec /bin/bash << EOT
  export HOME=/home/<%= user %>
  source /usr/local/share/chruby/chruby.sh
  chruby $FOREMAN_RUBY_VERSION
  exec <%= process.command %>
EOT
end script

